# garden_django/supervisord_djangoq.conf
[supervisord]
nodaemon=true
user=root

[program:django-queue-workers]
command=python manage.py qcluster
directory=.
autostart=true
autorestart=true
